# deploy/anti-ad-bot.service
[Unit]
Description=Anti-Ad Discord Bot
After=network.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/home/botuser/Anti-Ad
Environment="PATH=/home/botuser/Anti-Ad/venv/bin"
ExecStart=/home/botuser/Anti-Ad/venv/bin/python3 src/bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anti-ad-bot

[Install]
WantedBy=multi-user.target
